# docker/images/core/php7.3/conf/supervisor/conf.d/php.conf
[program:php-fpm]
command = /usr/sbin/php-fpm7 -F
user = root
autostart = true
